# networking/rc.local
#!/bin/sh

/sbin/iptables-restore < /etc/iptables.ipv4.nat

